# deb/scripts/prerm
#!/bin/bash

systemctl stop jjs-frontend
systemctl stop jjs-backend
systemctl disable jjs-frontend
systemctl disable jjs-backend
